# sample-config/systemd/esphere.socket
[Unit]
Description=ESP Sphere socket

[Socket]
ListenStream=/run/esphere.sock
SocketUser=www-data

[Install]
WantedBy=sockets.target
